# utility/files/etc/init.d/startup-service
#!/sbin/openrc-run

depend() {
  need net
  use logger
}
start() {
  ebegin "Starting startup script"
  /root/Scripts/startup &
  eend $?
}
stop() {
  ebegin "Startup script got exit signal."
  killall startup
  eend $?
}