# installer/files/etc/init.d/startup-service
#!/sbin/openrc-run

depend() {
  need docker net
  use logger
}
start() {
  ebegin "Starting startup script"
  /home/admin/Configs/Scripts/boot &
  eend $?
}
stop() {
  ebegin "Startup script stopped."
  eend $?
}
